# checkupdates.timer
[Timer]
OnCalendar=*-*-* 11:10
#Persistent=true

[Install]
WantedBy=timers.target
