# checkupdates.service
[Service]
ExecStart=/bin/checkupdates -d
